# config-files/systemd/user/parcimonie.service
[Unit]
Description=Parcimonie

[Service]
ExecStart=/home/kheops/apps/repository/git/parcimonie.sh/parcimonie.sh

[Install]
WantedBy=default.target
